# debian/postinst
#!/bin/bash
systemctl enable wolkabout-ELab && systemctl start wolkabout-ELab
