# debian/postinst
#!/bin/bash
systemctl enable wolkconnect-cpp && systemctl start wolkconnect-cpp
